# supervisor.conf
[unix_http_server]
file=/tmp/my-blog.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/my-blog.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/my-blog.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/my-blog.sock ; use a unix:// URL  for a unix socket

[program:my-blog1]
command= node index
;user=admin
autostart=true
autorestart=true
log_stdout=true
logfile_maxbytes=1MB
logfile_backups=10
stdout_logfile=./my-blog.log
logfile=./my-blog.log

# [program:new-wap2]
# command= node webserver.js -p 3009
# ;user=admin
# autostart=true
# autorestart=true
# log_stdout=true
# logfile_maxbytes=1MB
# logfile_backups=10
# stdout_logfile=./cout_3009.log
# logfile=./cout1_3009.log
